# systemd/bilalapp.be.service
[Unit]
Description=Bilal Backend Service
After=network.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=bilal
WorkingDirectory=/opt/bilalapp
Environment=PYTHONUNBUFFERED=1
Environment=BILAL_ENFORCE_SETTINGS_API=1
Environment=BILAL_RUN_SCHEDULER=0
ExecStart=/opt/bilalapp/.venv/bin/python3 /opt/bilalapp/server.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
