# systemd/bilalapp.fe.service
[Unit]
Description=Bilal Frontend Service (production static server)
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bilalapp
# Run the bundled Node static server directly (no npx) so the process runs as `bilal`.
# We recommend building artifacts before deployment and copying `dist/` to /opt/bilalapp.
# ExecStartPre removed to avoid rebuilding at service start; change if you prefer auto-build.
ExecStart=/usr/bin/node /opt/bilalapp/server-static.cjs
Restart=on-failure
RestartSec=5
User=bilal
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
